# assets/notes/com.alice-station.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alice-station.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mac/Desktop/Работа/alice-station-control/work/backend/venv/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mac/Desktop/Работа/alice-station-control/work/backend</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Ждём 10 сек после входа — чтобы Wi-Fi успел подключиться -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/mac/Library/Logs/alice-station.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mac/Library/Logs/alice-station.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mac/Desktop/Работа/alice-station-control/work/backend/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
